# alienware-m17xr3-ejectcd
#!/sbin/openrc-run

# This fixes the ejectcd button for Alienware M17xR3 laptops.
start() {
    ebegin "Starting Alienware M17xR3 fix"
    setkeycodes 89 161
}
